# SVAuth/IntegrationExamples/hotcrp/hotcrp/supervisord-hotcrp.conf
[program:hotcrp]
command=/start-hotcrp.sh
autostart=true
autorestart=false
